# net/nginx/root/etc/uci-defaults/nginx
#!/bin/sh

/etc/init.d/nginx disable
exit 0
